# deploy/prediction-pulse-drafter.service
[Unit]
Description=prediction-pulse drafter (writes post drafts on a cadence, dryrun)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Projects/prediction-pulse
# Cadence + per-cycle cap are the two knobs to tune (edit + daemon-reload + restart).
# WorkingDirectory makes .env (PULSE_PERSONA, ANTHROPIC_API_KEY), personas/, and the DB resolve.
ExecStart=/home/pi/Projects/prediction-pulse/.venv/bin/pulse draft --interval 3600 --limit 3
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
